# etc/supervisor/com.solordian.event.conf
[program:com.solardian.event]
command=php /data/projects/ldt/lvdiantong/app/bin/ptphp.php --cli=event --action=service
autostart=true
autorestart=true
startsecs=3